# trunk/experimental/gdm/.footprint
drwxr-xr-x	root/root	etc/
drwxr-xr-x	root/root	etc/dbus-1/
drwxr-xr-x	root/root	etc/dbus-1/system.d/
-rw-r--r--	root/root	etc/dbus-1/system.d/gdm.conf
drwxr-xr-x	root/root	etc/gconf/
drwxr-xr-x	root/root	etc/gconf/schemas/
-rw-r--r--	root/root	etc/gconf/schemas/gdm-simple-greeter.schemas
drwxr-xr-x	root/root	etc/gdm/
drwxr-xr-x	root/root	etc/gdm/Init/
-rwxr-xr-x	root/root	etc/gdm/Init/Default
drwxr-xr-x	root/root	etc/gdm/PostLogin/
-rwxr-xr-x	root/root	etc/gdm/PostLogin/Default.sample
drwxr-xr-x	root/root	etc/gdm/PostSession/
-rwxr-xr-x	root/root	etc/gdm/PostSession/Default
drwxr-xr-x	root/root	etc/gdm/PreSession/
-rwxr-xr-x	root/root	etc/gdm/PreSession/Default
-rwxr-xr-x	root/root	etc/gdm/Xsession
-rw-r--r--	root/root	etc/gdm/custom.conf
-rw-r--r--	root/root	etc/gdm/gdm.schemas
drwxr-xr-x	root/root	etc/pam.d/
-rw-r--r--	root/root	etc/pam.d/gdm
-rw-r--r--	root/root	etc/pam.d/gdm-autologin
drwxr-xr-x	root/root	etc/rc.d/
-rwxr-xr-x	root/root	etc/rc.d/gdm
drwxr-xr-x	root/root	usr/
drwxr-xr-x	root/root	usr/bin/
-rwxr-xr-x	root/root	usr/bin/gdm-screenshot
-rwxr-xr-x	root/root	usr/bin/gdmflexiserver
drwxr-xr-x	root/root	usr/lib/
drwxr-xr-x	root/root	usr/lib/bonobo/
drwxr-xr-x	root/root	usr/lib/bonobo/servers/
-rw-r--r--	root/root	usr/lib/bonobo/servers/GNOME_FastUserSwitchApplet.server
drwxr-xr-x	root/root	usr/libexec/
-rwxr-xr-x	root/root	usr/libexec/gdm-crash-logger
-rwxr-xr-x	root/root	usr/libexec/gdm-factory-slave
-rwxr-xr-x	root/root	usr/libexec/gdm-host-chooser
-rwxr-xr-x	root/root	usr/libexec/gdm-product-slave
-rwxr-xr-x	root/root	usr/libexec/gdm-session-worker
-rwxr-xr-x	root/root	usr/libexec/gdm-simple-chooser
-rwxr-xr-x	root/root	usr/libexec/gdm-simple-greeter
-rwxr-xr-x	root/root	usr/libexec/gdm-simple-slave
-rwxr-xr-x	root/root	usr/libexec/gdm-user-switch-applet
-rwxr-xr-x	root/root	usr/libexec/gdm-xdmcp-chooser-slave
drwxr-xr-x	root/root	usr/sbin/
-rwxr-xr-x	root/root	usr/sbin/gdm
-rwxr-xr-x	root/root	usr/sbin/gdm-binary
-rwxr-xr-x	root/root	usr/sbin/gdm-restart
-rwxr-xr-x	root/root	usr/sbin/gdm-safe-restart
-rwxr-xr-x	root/root	usr/sbin/gdm-stop
drwxr-xr-x	root/root	usr/share/
drwxr-xr-x	root/root	usr/share/gdm/
drwxr-xr-x	root/root	usr/share/gdm/autostart/
drwxr-xr-x	root/root	usr/share/gdm/autostart/LoginWindow/
-rw-r--r--	root/root	usr/share/gdm/autostart/LoginWindow/at-spi-registryd-wrapper.desktop
-rw-r--r--	root/root	usr/share/gdm/autostart/LoginWindow/gdm-simple-greeter.desktop
-rw-r--r--	root/root	usr/share/gdm/autostart/LoginWindow/gnome-mag.desktop
-rw-r--r--	root/root	usr/share/gdm/autostart/LoginWindow/gnome-power-manager.desktop
-rw-r--r--	root/root	usr/share/gdm/autostart/LoginWindow/gnome-settings-daemon.desktop
-rw-r--r--	root/root	usr/share/gdm/autostart/LoginWindow/gok.desktop
-rw-r--r--	root/root	usr/share/gdm/autostart/LoginWindow/metacity.desktop
-rw-r--r--	root/root	usr/share/gdm/autostart/LoginWindow/orca-screen-reader.desktop
-rw-r--r--	root/root	usr/share/gdm/gdb-cmd
-rw-r--r--	root/root	usr/share/gdm/gdm-greeter-login-window.glade
drwxr-xr-x	root/root	usr/share/gdm/themes/
-rw-r--r--	root/root	usr/share/gdm/themes/GdmGreeterTheme.desktop
-rw-r--r--	root/root	usr/share/gdm/themes/crux.png
-rw-r--r--	root/root	usr/share/gdm/themes/crux.xml
-rw-r--r--	root/root	usr/share/gdm/themes/options.png
-rw-r--r--	root/root	usr/share/gdm/themes/screenshot.jpg
drwxr-xr-x	root/root	usr/share/gnome-2.0/
drwxr-xr-x	root/root	usr/share/gnome-2.0/ui/
-rw-r--r--	root/root	usr/share/gnome-2.0/ui/GNOME_FastUserSwitchApplet.xml
drwxr-xr-x	root/root	usr/share/gnome/
drwxr-xr-x	root/root	usr/share/gnome/help/
drwxr-xr-x	root/root	usr/share/gnome/help/gdm/
drwxr-xr-x	root/root	usr/share/gnome/help/gdm/C/
-rw-r--r--	root/root	usr/share/gnome/help/gdm/C/gdm.xml
-rw-r--r--	root/root	usr/share/gnome/help/gdm/C/legal.xml
drwxr-xr-x	root/root	usr/share/gnome/help/gdm/de/
-rw-r--r--	root/root	usr/share/gnome/help/gdm/de/gdm.xml
drwxr-xr-x	root/root	usr/share/gnome/help/gdm/en_GB/
-rw-r--r--	root/root	usr/share/gnome/help/gdm/en_GB/gdm.xml
drwxr-xr-x	root/root	usr/share/gnome/help/gdm/es/
-rw-r--r--	root/root	usr/share/gnome/help/gdm/es/gdm.xml
drwxr-xr-x	root/root	usr/share/gnome/help/gdm/fr/
-rw-r--r--	root/root	usr/share/gnome/help/gdm/fr/gdm.xml
drwxr-xr-x	root/root	usr/share/gnome/help/gdm/ko/
-rw-r--r--	root/root	usr/share/gnome/help/gdm/ko/gdm.xml
drwxr-xr-x	root/root	usr/share/gnome/help/gdm/oc/
-rw-r--r--	root/root	usr/share/gnome/help/gdm/oc/gdm.xml
drwxr-xr-x	root/root	usr/share/gnome/help/gdm/ru/
-rw-r--r--	root/root	usr/share/gnome/help/gdm/ru/gdm.xml
drwxr-xr-x	root/root	usr/share/gnome/help/gdm/sv/
-rw-r--r--	root/root	usr/share/gnome/help/gdm/sv/gdm.xml
drwxr-xr-x	root/root	usr/share/gnome/help/gdm/uk/
-rw-r--r--	root/root	usr/share/gnome/help/gdm/uk/gdm.xml
drwxr-xr-x	root/root	usr/share/icons/
drwxr-xr-x	root/root	usr/share/icons/hicolor/
drwxr-xr-x	root/root	usr/share/icons/hicolor/16x16/
drwxr-xr-x	root/root	usr/share/icons/hicolor/16x16/apps/
-rw-r--r--	root/root	usr/share/icons/hicolor/16x16/apps/gdm-xnest.png
drwxr-xr-x	root/root	usr/share/icons/hicolor/32x32/
drwxr-xr-x	root/root	usr/share/icons/hicolor/32x32/apps/
-rw-r--r--	root/root	usr/share/icons/hicolor/32x32/apps/gdm-setup.png
-rw-r--r--	root/root	usr/share/icons/hicolor/32x32/apps/gdm-xnest.png
drwxr-xr-x	root/root	usr/share/pixmaps/
drwxr-xr-x	root/root	usr/share/pixmaps/faces/
-rw-r--r--	root/root	usr/share/pixmaps/faces/astronaut.jpg
-rw-r--r--	root/root	usr/share/pixmaps/faces/baseball.png
-rw-r--r--	root/root	usr/share/pixmaps/faces/butterfly.png
-rw-r--r--	root/root	usr/share/pixmaps/faces/cat-eye.jpg
-rw-r--r--	root/root	usr/share/pixmaps/faces/chess.jpg
-rw-r--r--	root/root	usr/share/pixmaps/faces/coffee.jpg
-rw-r--r--	root/root	usr/share/pixmaps/faces/dice.jpg
-rw-r--r--	root/root	usr/share/pixmaps/faces/energy-arc.jpg
-rw-r--r--	root/root	usr/share/pixmaps/faces/fish.jpg
-rw-r--r--	root/root	usr/share/pixmaps/faces/flake.jpg
-rw-r--r--	root/root	usr/share/pixmaps/faces/flower.jpg
-rw-r--r--	root/root	usr/share/pixmaps/faces/grapes.jpg
-rw-r--r--	root/root	usr/share/pixmaps/faces/guitar.jpg
-rw-r--r--	root/root	usr/share/pixmaps/faces/launch.jpg
-rw-r--r--	root/root	usr/share/pixmaps/faces/leaf.jpg
-rw-r--r--	root/root	usr/share/pixmaps/faces/lightning.jpg
-rw-r--r--	root/root	usr/share/pixmaps/faces/penguin.jpg
-rw-r--r--	root/root	usr/share/pixmaps/faces/puppy.jpg
-rw-r--r--	root/root	usr/share/pixmaps/faces/sky.jpg
-rw-r--r--	root/root	usr/share/pixmaps/faces/soccerball.png
-rw-r--r--	root/root	usr/share/pixmaps/faces/sunflower.jpg
-rw-r--r--	root/root	usr/share/pixmaps/faces/sunset.jpg
-rw-r--r--	root/root	usr/share/pixmaps/faces/tennis-ball.png
-rw-r--r--	root/root	usr/share/pixmaps/faces/yellow-rose.jpg
-rw-r--r--	root/root	usr/share/pixmaps/gdm-foot-logo.png
-rw-r--r--	root/root	usr/share/pixmaps/gdm-setup.png
-rw-r--r--	root/root	usr/share/pixmaps/gdm-xnest.png
-rw-r--r--	root/root	usr/share/pixmaps/gdm.png
-rw-r--r--	root/root	usr/share/pixmaps/nobody.png
-rw-r--r--	root/root	usr/share/pixmaps/nohost.png
drwxr-xr-x	root/root	var/
drwxrwx--T	root/gdm	var/gdm/
drwxr-xr-x	root/root	var/lib/
drwxrwx--T	root/gdm	var/lib/gdm/
drwxr-x--T	root/gdm	var/lib/gdm/.gconf.mandatory/
-rw-r----T	root/gdm	var/lib/gdm/.gconf.mandatory/%gconf-tree.xml
-rw-r--r--	root/root	var/lib/gdm/.gconf.path
drwxr-xr-x	root/root	var/log/
drwxr-xr-x	root/root	var/log/gdm/
drwxr-xr-x	root/root	var/run/
drwxrwxrwt	root/gdm	var/run/gdm/
