# experimental/bluez-utils/.footprint
drwxr-xr-x	root/root	etc/
drwxr-xr-x	root/root	etc/bluetooth/
-rw-r--r--	root/root	etc/bluetooth/hcid.conf
-rw-r--r--	root/root	etc/bluetooth/rfcomm.conf
drwxr-xr-x	root/root	etc/rc.d/
-rwxr-xr-x	root/root	etc/rc.d/bluetooth
drwxr-xr-x	root/root	usr/
drwxr-xr-x	root/root	usr/bin/
-rwxr-xr-x	root/root	usr/bin/ciptool
-rwxr-xr-x	root/root	usr/bin/dund
-rwxr-xr-x	root/root	usr/bin/hcitool
-rwxr-xr-x	root/root	usr/bin/hidd
-rwxr-xr-x	root/root	usr/bin/l2ping
-rwxr-xr-x	root/root	usr/bin/pand
-rwxr-xr-x	root/root	usr/bin/rfcomm
-rwxr-xr-x	root/root	usr/bin/sdptool
drwxr-xr-x	root/root	usr/etc/
drwxr-xr-x	root/root	usr/etc/dbus-1/
drwxr-xr-x	root/root	usr/etc/dbus-1/system.d/
-rw-r--r--	root/root	usr/etc/dbus-1/system.d/bluetooth.conf
drwxr-xr-x	root/root	usr/lib/
drwxr-xr-x	root/root	usr/lib/alsa-lib/
-rwxr-xr-x	root/root	usr/lib/alsa-lib/libasound_module_ctl_bluetooth.la
-rwxr-xr-x	root/root	usr/lib/alsa-lib/libasound_module_ctl_bluetooth.so
-rwxr-xr-x	root/root	usr/lib/alsa-lib/libasound_module_pcm_bluetooth.la
-rwxr-xr-x	root/root	usr/lib/alsa-lib/libasound_module_pcm_bluetooth.so
drwxr-xr-x	root/root	usr/lib/bluetooth/
drwxr-xr-x	root/root	usr/lib/bluetooth/plugins/
-rwxr-xr-x	root/root	usr/lib/bluetooth/plugins/audio.la
-rwxr-xr-x	root/root	usr/lib/bluetooth/plugins/audio.so
-rwxr-xr-x	root/root	usr/lib/bluetooth/plugins/input.la
-rwxr-xr-x	root/root	usr/lib/bluetooth/plugins/input.so
-rwxr-xr-x	root/root	usr/lib/bluetooth/plugins/network.la
-rwxr-xr-x	root/root	usr/lib/bluetooth/plugins/network.so
-rwxr-xr-x	root/root	usr/lib/bluetooth/plugins/serial.la
-rwxr-xr-x	root/root	usr/lib/bluetooth/plugins/serial.so
drwxr-xr-x	root/root	usr/lib/cups/
drwxr-xr-x	root/root	usr/lib/cups/backend/
-rwxr-xr-x	root/root	usr/lib/cups/backend/bluetooth
drwxr-xr-x	root/root	usr/man/
drwxr-xr-x	root/root	usr/man/man1/
-rw-r--r--	root/root	usr/man/man1/ciptool.1.gz
-rw-r--r--	root/root	usr/man/man1/dund.1.gz
-rw-r--r--	root/root	usr/man/man1/hcitool.1.gz
-rw-r--r--	root/root	usr/man/man1/hidd.1.gz
-rw-r--r--	root/root	usr/man/man1/l2ping.1.gz
-rw-r--r--	root/root	usr/man/man1/pand.1.gz
-rw-r--r--	root/root	usr/man/man1/rfcomm.1.gz
-rw-r--r--	root/root	usr/man/man1/sdptool.1.gz
drwxr-xr-x	root/root	usr/man/man5/
-rw-r--r--	root/root	usr/man/man5/hcid.conf.5.gz
drwxr-xr-x	root/root	usr/man/man8/
-rw-r--r--	root/root	usr/man/man8/hciattach.8.gz
-rw-r--r--	root/root	usr/man/man8/hciconfig.8.gz
-rw-r--r--	root/root	usr/man/man8/hcid.8.gz
-rw-r--r--	root/root	usr/man/man8/hid2hci.8.gz
drwxr-xr-x	root/root	usr/sbin/
-rwxr-xr-x	root/root	usr/sbin/hciattach
-rwxr-xr-x	root/root	usr/sbin/hciconfig
-rwxr-xr-x	root/root	usr/sbin/hcid
-rwxr-xr-x	root/root	usr/sbin/hid2hci
drwxr-xr-x	root/root	var/
drwxr-xr-x	root/root	var/lib/
drwxr-xr-x	root/root	var/lib/bluetooth/
