# trunk/ports/acpi-support/.footprint
drwxr-xr-x	root/root	etc/
drwxr-xr-x	root/root	etc/acpi/
drwxr-xr-x	root/root	etc/acpi/ac.d/
-rwxr-xr-x	root/root	etc/acpi/ac.d/90-hdparm.sh
-rwxr-xr-x	root/root	etc/acpi/always-mute.sh
-rwxr-xr-x	root/root	etc/acpi/asus-brn-down.sh
-rwxr-xr-x	root/root	etc/acpi/asus-brn-up.sh
-rwxr-xr-x	root/root	etc/acpi/asus-touchpad.sh
-rwxr-xr-x	root/root	etc/acpi/asus-wireless-2.sh
-rwxr-xr-x	root/root	etc/acpi/asus-wireless.sh
drwxr-xr-x	root/root	etc/acpi/battery.d/
-rwxr-xr-x	root/root	etc/acpi/battery.d/90-hdparm.sh
-rwxr-xr-x	root/root	etc/acpi/batterybtn.sh
-rwxr-xr-x	root/root	etc/acpi/ejectbtn.sh
drwxr-xr-x	root/root	etc/acpi/events/
-rw-r--r--	root/root	etc/acpi/events/ac
-rw-r--r--	root/root	etc/acpi/events/asus-a6u-touchpad
-rw-r--r--	root/root	etc/acpi/events/asus-brightness-down
-rw-r--r--	root/root	etc/acpi/events/asus-brightness-up
-rw-r--r--	root/root	etc/acpi/events/asus-eee-volume-down
-rw-r--r--	root/root	etc/acpi/events/asus-eee-volume-mute
-rw-r--r--	root/root	etc/acpi/events/asus-eee-volume-up
-rw-r--r--	root/root	etc/acpi/events/asus-internet
-rw-r--r--	root/root	etc/acpi/events/asus-lock
-rw-r--r--	root/root	etc/acpi/events/asus-mail
-rw-r--r--	root/root	etc/acpi/events/asus-media-eject
-rw-r--r--	root/root	etc/acpi/events/asus-media-next
-rw-r--r--	root/root	etc/acpi/events/asus-media-play-pause
-rw-r--r--	root/root	etc/acpi/events/asus-media-prev
-rw-r--r--	root/root	etc/acpi/events/asus-media-stop
-rw-r--r--	root/root	etc/acpi/events/asus-rotate
-rw-r--r--	root/root	etc/acpi/events/asus-touchpad
-rw-r--r--	root/root	etc/acpi/events/asus-video
-rw-r--r--	root/root	etc/acpi/events/asus-volume-down
-rw-r--r--	root/root	etc/acpi/events/asus-volume-mute
-rw-r--r--	root/root	etc/acpi/events/asus-volume-up
-rw-r--r--	root/root	etc/acpi/events/asus-wireless
-rw-r--r--	root/root	etc/acpi/events/asus-wireless-off
-rw-r--r--	root/root	etc/acpi/events/asus-wireless-on
-rw-r--r--	root/root	etc/acpi/events/battery
-rw-r--r--	root/root	etc/acpi/events/ibm-hibernatebtn
-rw-r--r--	root/root	etc/acpi/events/ibm-lockbtn
-rw-r--r--	root/root	etc/acpi/events/ibm-sleepbtn
-rw-r--r--	root/root	etc/acpi/events/ibm-videobtn
-rw-r--r--	root/root	etc/acpi/events/ibm-wireless
-rw-r--r--	root/root	etc/acpi/events/lenovo-lockbtn
-rw-r--r--	root/root	etc/acpi/events/lenovo-touchpad
-rw-r--r--	root/root	etc/acpi/events/lidbtn
-rw-r--r--	root/root	etc/acpi/events/panasonic-brightness-down
-rw-r--r--	root/root	etc/acpi/events/panasonic-brightness-up
-rw-r--r--	root/root	etc/acpi/events/panasonic-hibernatebtn
-rw-r--r--	root/root	etc/acpi/events/panasonic-lockbtn
-rw-r--r--	root/root	etc/acpi/events/panasonic-mute
-rw-r--r--	root/root	etc/acpi/events/panasonic-sleepbtn
-rw-r--r--	root/root	etc/acpi/events/panasonic-volume-down
-rw-r--r--	root/root	etc/acpi/events/panasonic-volume-up
-rw-r--r--	root/root	etc/acpi/events/sleepbtn
-rw-r--r--	root/root	etc/acpi/events/sony-brightness-down
-rw-r--r--	root/root	etc/acpi/events/sony-brightness-up
-rw-r--r--	root/root	etc/acpi/events/sony-eject
-rw-r--r--	root/root	etc/acpi/events/sony-hibernate
-rw-r--r--	root/root	etc/acpi/events/sony-mute
-rw-r--r--	root/root	etc/acpi/events/sony-sleep
-rw-r--r--	root/root	etc/acpi/events/sony-volume-down
-rw-r--r--	root/root	etc/acpi/events/sony-volume-up
-rw-r--r--	root/root	etc/acpi/events/thinkpad-brightness-down
-rw-r--r--	root/root	etc/acpi/events/thinkpad-brightness-up
-rw-r--r--	root/root	etc/acpi/events/thinkpad-cmos
-rw-r--r--	root/root	etc/acpi/events/thinkpad-mute
-rw-r--r--	root/root	etc/acpi/events/thinkpad-thinklight
-rw-r--r--	root/root	etc/acpi/events/thinkpad-thinkpad
-rw-r--r--	root/root	etc/acpi/events/thinkpad-volume-down
-rw-r--r--	root/root	etc/acpi/events/thinkpad-volume-up
-rw-r--r--	root/root	etc/acpi/events/thinkpad-zoom
-rw-r--r--	root/root	etc/acpi/events/tosh-battery
-rw-r--r--	root/root	etc/acpi/events/tosh-brightness-down
-rw-r--r--	root/root	etc/acpi/events/tosh-brightness-up
-rw-r--r--	root/root	etc/acpi/events/tosh-hibernate
-rw-r--r--	root/root	etc/acpi/events/tosh-ibutton
-rw-r--r--	root/root	etc/acpi/events/tosh-lock
-rw-r--r--	root/root	etc/acpi/events/tosh-mail
-rw-r--r--	root/root	etc/acpi/events/tosh-media
-rw-r--r--	root/root	etc/acpi/events/tosh-mute
-rw-r--r--	root/root	etc/acpi/events/tosh-next
-rw-r--r--	root/root	etc/acpi/events/tosh-play
-rw-r--r--	root/root	etc/acpi/events/tosh-prev
-rw-r--r--	root/root	etc/acpi/events/tosh-sleep
-rw-r--r--	root/root	etc/acpi/events/tosh-stop
-rw-r--r--	root/root	etc/acpi/events/tosh-wireless
-rw-r--r--	root/root	etc/acpi/events/tosh-www
-rw-r--r--	root/root	etc/acpi/events/video_brightnessdown
-rw-r--r--	root/root	etc/acpi/events/video_brightnessup
-rw-r--r--	root/root	etc/acpi/events/videobtn
-rwxr-xr-x	root/root	etc/acpi/hibernate.sh
-rwxr-xr-x	root/root	etc/acpi/hibernatebtn.sh
-rwxr-xr-x	root/root	etc/acpi/ibm-wireless.sh
-rwxr-xr-x	root/root	etc/acpi/if-asus-eee.sh
-rwxr-xr-x	root/root	etc/acpi/lid.sh
-rwxr-xr-x	root/root	etc/acpi/lockbtn.sh
-rwxr-xr-x	root/root	etc/acpi/mailbtn.sh
-rwxr-xr-x	root/root	etc/acpi/mediabtn.sh
-rwxr-xr-x	root/root	etc/acpi/mutebtn.sh
-rwxr-xr-x	root/root	etc/acpi/nextbtn.sh
-rwxr-xr-x	root/root	etc/acpi/panabright.sh
-rwxr-xr-x	root/root	etc/acpi/panapower.sh
-rwxr-xr-x	root/root	etc/acpi/playbtn.sh
-rwxr-xr-x	root/root	etc/acpi/power.sh
-rwxr-xr-x	root/root	etc/acpi/prepare.sh
-rwxr-xr-x	root/root	etc/acpi/prevbtn.sh
drwxr-xr-x	root/root	etc/acpi/resume.d/
-rwxr-xr-x	root/root	etc/acpi/resume.d/10-thinkpad-standby-led.sh
-rwxr-xr-x	root/root	etc/acpi/resume.d/11-video-pci-restore.sh
-rwxr-xr-x	root/root	etc/acpi/resume.d/15-video-post.sh
-rwxr-xr-x	root/root	etc/acpi/resume.d/17-video-restore.sh
-rwxr-xr-x	root/root	etc/acpi/resume.d/35-modules-load.sh
-rwxr-xr-x	root/root	etc/acpi/resume.d/40-infra-red.sh
-rwxr-xr-x	root/root	etc/acpi/resume.d/50-framebuffer-enable.sh
-rwxr-xr-x	root/root	etc/acpi/resume.d/50-time.sh
-rwxr-xr-x	root/root	etc/acpi/resume.d/50-tosh-restore-brightness.sh
-rwxr-xr-x	root/root	etc/acpi/resume.d/55-screen.sh
-rwxr-xr-x	root/root	etc/acpi/resume.d/58-proc-sysfs-restore-state.sh
-rwxr-xr-x	root/root	etc/acpi/resume.d/60-asus-wireless-led.sh
-rwxr-xr-x	root/root	etc/acpi/resume.d/62-ifup.sh
-rwxr-xr-x	root/root	etc/acpi/resume.d/65-console.sh
-rwxr-xr-x	root/root	etc/acpi/resume.d/67-sound.sh
-rwxr-xr-x	root/root	etc/acpi/resume.d/69-services.sh
-rwxr-xr-x	root/root	etc/acpi/resume.d/72-acpi-pain.sh
-rwxr-xr-x	root/root	etc/acpi/resume.d/89-asus-restore-brightness.sh
-rwxr-xr-x	root/root	etc/acpi/resume.d/90-hdparm.sh
-rwxr-xr-x	root/root	etc/acpi/resume.d/90-thinkpad-unstandby-led.sh
-rwxr-xr-x	root/root	etc/acpi/resume.d/90-xscreensaver.sh
-rwxr-xr-x	root/root	etc/acpi/resume.d/95-laptop-mode.sh
-rwxr-xr-x	root/root	etc/acpi/resume.d/98-acpi-unlock.sh
-rwxr-xr-x	root/root	etc/acpi/resume.sh
-rwxr-xr-x	root/root	etc/acpi/rotatescreen.sh
-rwxr-xr-x	root/root	etc/acpi/screenblank.sh
-rwxr-xr-x	root/root	etc/acpi/sleep.sh
-rwxr-xr-x	root/root	etc/acpi/sleepbtn.sh
-rwxr-xr-x	root/root	etc/acpi/sonybright.sh
drwxr-xr-x	root/root	etc/acpi/start.d/
-rwxr-xr-x	root/root	etc/acpi/start.d/10-save-dmidecode.sh
-rwxr-xr-x	root/root	etc/acpi/start.d/60-asus-wireless-led.sh
-rwxr-xr-x	root/root	etc/acpi/start.d/90-hdparm.sh
-rwxr-xr-x	root/root	etc/acpi/stopbtn.sh
drwxr-xr-x	root/root	etc/acpi/suspend.d/
-rwxr-xr-x	root/root	etc/acpi/suspend.d/01-laptop-mode.sh
-rwxr-xr-x	root/root	etc/acpi/suspend.d/05-acpi-lock.sh
-rwxr-xr-x	root/root	etc/acpi/suspend.d/10-thinkpad-standby-led.sh
-rwxr-xr-x	root/root	etc/acpi/suspend.d/30-proc-sysfs-save-state.sh
-rwxr-xr-x	root/root	etc/acpi/suspend.d/50-irda-stop.sh
-rwxr-xr-x	root/root	etc/acpi/suspend.d/50-time.sh
-rwxr-xr-x	root/root	etc/acpi/suspend.d/50-tosh-save-brightness.sh
-rwxr-xr-x	root/root	etc/acpi/suspend.d/55-down-interfaces.sh
-rwxr-xr-x	root/root	etc/acpi/suspend.d/60-generate-modules-list.sh
-rwxr-xr-x	root/root	etc/acpi/suspend.d/65-services-stop.sh
-rwxr-xr-x	root/root	etc/acpi/suspend.d/70-modules-unload.sh
-rwxr-xr-x	root/root	etc/acpi/suspend.d/75-console-switch.sh
-rwxr-xr-x	root/root	etc/acpi/suspend.d/80-video-vesa-state.sh
-rwxr-xr-x	root/root	etc/acpi/suspend.d/81-video-pci-state.sh
-rwxr-xr-x	root/root	etc/acpi/suspend.d/85-alsa-state.sh
-rwxr-xr-x	root/root	etc/acpi/suspend.d/90-framebuffer-stop.sh
-rwxr-xr-x	root/root	etc/acpi/thinkpad-brightness-down.sh
-rwxr-xr-x	root/root	etc/acpi/thinkpad-brightness-up.sh
-rwxr-xr-x	root/root	etc/acpi/thinkpad-lockorbattery.sh
-rwxr-xr-x	root/root	etc/acpi/thinkpad-stretchortouchpad.sh
-rwxr-xr-x	root/root	etc/acpi/thinkpad-thinklight.sh
-rwxr-xr-x	root/root	etc/acpi/thinkpad-thinkpad.sh
-rwxr-xr-x	root/root	etc/acpi/thinkpad-zoom.sh
-rwxr-xr-x	root/root	etc/acpi/tosh-wireless.sh
-rwxr-xr-x	root/root	etc/acpi/toshbright.sh
-rwxr-xr-x	root/root	etc/acpi/video_brightnessdown.sh
-rwxr-xr-x	root/root	etc/acpi/video_brightnessup.sh
-rwxr-xr-x	root/root	etc/acpi/videobtn.sh
-rwxr-xr-x	root/root	etc/acpi/voldownbtn.sh
-rwxr-xr-x	root/root	etc/acpi/volupbtn.sh
-rwxr-xr-x	root/root	etc/acpi/webbtn.sh
-rwxr-xr-x	root/root	etc/acpi/wireless.sh
drwxr-xr-x	root/root	etc/default/
-rw-r--r--	root/root	etc/default/acpi-support
drwxr-xr-x	root/root	etc/modprobe.d/
-rw-r--r--	root/root	etc/modprobe.d/thinkpad_acpi.modprobe
drwxr-xr-x	root/root	etc/rc.d/
-rwxr-xr-x	root/root	etc/rc.d/vbesave
drwxr-xr-x	root/root	usr/
drwxr-xr-x	root/root	usr/bin/
-rwxr-xr-x	root/root	usr/bin/acpi_fakekey
drwxr-xr-x	root/root	usr/share/
drwxr-xr-x	root/root	usr/share/acpi-support/
-rw-r--r--	root/root	usr/share/acpi-support/56548.config
-rw-r--r--	root/root	usr/share/acpi-support/ASUSTeK Computer Inc..config
-rw-r--r--	root/root	usr/share/acpi-support/Acer, inc..config
-rw-r--r--	root/root	usr/share/acpi-support/Acer.config
-rw-r--r--	root/root	usr/share/acpi-support/Compaq.config
-rw-r--r--	root/root	usr/share/acpi-support/Dell Computer Corporation.config
-rw-r--r--	root/root	usr/share/acpi-support/Dell Inc..config
-rw-r--r--	root/root	usr/share/acpi-support/ECS.config
-rw-r--r--	root/root	usr/share/acpi-support/FUJITSU SIEMENS.config
-rw-r--r--	root/root	usr/share/acpi-support/FUJITSU.config
-rw-r--r--	root/root	usr/share/acpi-support/Hewlett-Packard.config
-rw-r--r--	root/root	usr/share/acpi-support/IBM.config
-rw-r--r--	root/root	usr/share/acpi-support/NEC Computers International.config
-rw-r--r--	root/root	usr/share/acpi-support/SHARP Corporation.config
-rw-r--r--	root/root	usr/share/acpi-support/Samsung Electronics.config
-rw-r--r--	root/root	usr/share/acpi-support/Sony Corporation.config
-rw-r--r--	root/root	usr/share/acpi-support/TOSHIBA.config
-rw-r--r--	root/root	usr/share/acpi-support/device-funcs
-rw-r--r--	root/root	usr/share/acpi-support/key-constants
-rw-r--r--	root/root	usr/share/acpi-support/screenblank
-rw-r--r--	root/root	usr/share/acpi-support/state-funcs
-rwxr-xr-x	root/root	usr/share/acpi-support/suspendorhibernate
