# trunk/ports/dictd/.footprint
drwxr-xr-x	root/root	etc/
drwxr-xr-x	root/root	etc/dict/
-rwxr-xr-x	root/root	etc/dict/dict.conf
-rwxr-xr-x	root/root	etc/dict/dictd.conf
drwxr-xr-x	root/root	etc/rc.d/
-rwxr-xr-x	root/root	etc/rc.d/dictd
drwxr-xr-x	root/root	usr/
drwxr-xr-x	root/root	usr/bin/
-rwxr-xr-x	root/root	usr/bin/colorit
-rwxr-xr-x	root/root	usr/bin/dict
-rwxr-xr-x	root/root	usr/bin/dictdplugin-config
-rwxr-xr-x	root/root	usr/bin/dictfmt
-rwxr-xr-x	root/root	usr/bin/dictfmt_index2suffix
-rwxr-xr-x	root/root	usr/bin/dictfmt_index2word
-rwxr-xr-x	root/root	usr/bin/dictl
-rwxr-xr-x	root/root	usr/bin/dictunformat
-rwxr-xr-x	root/root	usr/bin/dictzip
drwxr-xr-x	root/root	usr/include/
-rw-r--r--	root/root	usr/include/dictdplugin.h
drwxr-xr-x	root/root	usr/lib/
drwxr-xr-x	root/root	usr/man/
drwxr-xr-x	root/root	usr/man/man1/
-rw-r--r--	root/root	usr/man/man1/colorit.1.gz
-rw-r--r--	root/root	usr/man/man1/dict.1.gz
-rw-r--r--	root/root	usr/man/man1/dictfmt.1.gz
-rw-r--r--	root/root	usr/man/man1/dictfmt_index2suffix.1.gz
-rw-r--r--	root/root	usr/man/man1/dictfmt_index2word.1.gz
-rw-r--r--	root/root	usr/man/man1/dictl.1.gz
-rw-r--r--	root/root	usr/man/man1/dictunformat.1.gz
-rw-r--r--	root/root	usr/man/man1/dictzip.1.gz
drwxr-xr-x	root/root	usr/man/man8/
-rw-r--r--	root/root	usr/man/man8/dictd.8.gz
drwxr-xr-x	root/root	usr/sbin/
-rwxr-xr-x	root/root	usr/sbin/dictd
