# trunk/ports/ferm/.footprint
drwxr-xr-x	root/root	etc/
drwxr-xr-x	root/root	etc/ferm/
-rw-r--r--	root/root	etc/ferm/ferm.conf
drwxr-xr-x	root/root	etc/rc.d/
-rwxr-xr-x	root/root	etc/rc.d/ferm
drwxr-xr-x	root/root	usr/
drwxr-xr-x	root/root	usr/man/
drwxr-xr-x	root/root	usr/man/man1/
-rw-r--r--	root/root	usr/man/man1/ferm.1.gz
-rw-r--r--	root/root	usr/man/man1/import-ferm.1.gz
drwxr-xr-x	root/root	usr/sbin/
-rwxr-xr-x	root/root	usr/sbin/ferm
-rwxr-xr-x	root/root	usr/sbin/import-ferm
