# trunk/ports/empathy/.footprint
drwxr-xr-x	root/root	etc/
drwxr-xr-x	root/root	etc/bluetooth/
-rw-r--r--	root/root	etc/bluetooth/main.conf
-rw-r--r--	root/root	etc/bluetooth/rfcomm.conf
drwxr-xr-x	root/root	etc/dbus-1/
drwxr-xr-x	root/root	etc/dbus-1/system.d/
-rw-r--r--	root/root	etc/dbus-1/system.d/bluetooth.conf
drwxr-xr-x	root/root	etc/rc.d/
-rwxr-xr-x	root/root	etc/rc.d/bluez
drwxr-xr-x	root/root	usr/
drwxr-xr-x	root/root	usr/bin/
-rwxr-xr-x	root/root	usr/bin/ciptool
-rwxr-xr-x	root/root	usr/bin/dund
-rwxr-xr-x	root/root	usr/bin/hcitool
-rwxr-xr-x	root/root	usr/bin/hidd
-rwxr-xr-x	root/root	usr/bin/l2ping
-rwxr-xr-x	root/root	usr/bin/pand
-rwxr-xr-x	root/root	usr/bin/rfcomm
-rwxr-xr-x	root/root	usr/bin/sdptool
drwxr-xr-x	root/root	usr/include/
drwxr-xr-x	root/root	usr/include/bluetooth/
-rw-r--r--	root/root	usr/include/bluetooth/bluetooth.h
-rw-r--r--	root/root	usr/include/bluetooth/bnep.h
-rw-r--r--	root/root	usr/include/bluetooth/cmtp.h
-rw-r--r--	root/root	usr/include/bluetooth/hci.h
-rw-r--r--	root/root	usr/include/bluetooth/hci_lib.h
-rw-r--r--	root/root	usr/include/bluetooth/hidp.h
-rw-r--r--	root/root	usr/include/bluetooth/l2cap.h
-rw-r--r--	root/root	usr/include/bluetooth/rfcomm.h
-rw-r--r--	root/root	usr/include/bluetooth/sco.h
-rw-r--r--	root/root	usr/include/bluetooth/sdp.h
-rw-r--r--	root/root	usr/include/bluetooth/sdp_lib.h
drwxr-xr-x	root/root	usr/lib/
drwxr-xr-x	root/root	usr/lib/alsa-lib/
-rwxr-xr-x	root/root	usr/lib/alsa-lib/libasound_module_ctl_bluetooth.la
-rwxr-xr-x	root/root	usr/lib/alsa-lib/libasound_module_ctl_bluetooth.so
-rwxr-xr-x	root/root	usr/lib/alsa-lib/libasound_module_pcm_bluetooth.la
-rwxr-xr-x	root/root	usr/lib/alsa-lib/libasound_module_pcm_bluetooth.so
drwxr-xr-x	root/root	usr/lib/bluetooth/
drwxr-xr-x	root/root	usr/lib/bluetooth/plugins/
-rwxr-xr-x	root/root	usr/lib/bluetooth/plugins/audio.la
-rwxr-xr-x	root/root	usr/lib/bluetooth/plugins/audio.so
-rwxr-xr-x	root/root	usr/lib/bluetooth/plugins/hal.la
-rwxr-xr-x	root/root	usr/lib/bluetooth/plugins/hal.so
-rwxr-xr-x	root/root	usr/lib/bluetooth/plugins/input.la
-rwxr-xr-x	root/root	usr/lib/bluetooth/plugins/input.so
-rwxr-xr-x	root/root	usr/lib/bluetooth/plugins/network.la
-rwxr-xr-x	root/root	usr/lib/bluetooth/plugins/network.so
-rwxr-xr-x	root/root	usr/lib/bluetooth/plugins/serial.la
-rwxr-xr-x	root/root	usr/lib/bluetooth/plugins/serial.so
-rwxr-xr-x	root/root	usr/lib/bluetooth/plugins/service.la
-rwxr-xr-x	root/root	usr/lib/bluetooth/plugins/service.so
drwxr-xr-x	root/root	usr/lib/cups/
drwxr-xr-x	root/root	usr/lib/cups/backend/
-rwxr-xr-x	root/root	usr/lib/cups/backend/bluetooth
drwxr-xr-x	root/root	usr/lib/gstreamer-0.10/
-rwxr-xr-x	root/root	usr/lib/gstreamer-0.10/libgstbluetooth.la
-rwxr-xr-x	root/root	usr/lib/gstreamer-0.10/libgstbluetooth.so
-rwxr-xr-x	root/root	usr/lib/libbluetooth.la
lrwxrwxrwx	root/root	usr/lib/libbluetooth.so -> libbluetooth.so.3.2.1
lrwxrwxrwx	root/root	usr/lib/libbluetooth.so.3 -> libbluetooth.so.3.2.1
-rwxr-xr-x	root/root	usr/lib/libbluetooth.so.3.2.1
drwxr-xr-x	root/root	usr/lib/pkgconfig/
-rw-r--r--	root/root	usr/lib/pkgconfig/bluez.pc
drwxr-xr-x	root/root	usr/man1/
-rw-r--r--	root/root	usr/man1/ciptool.1
-rw-r--r--	root/root	usr/man1/dund.1
-rw-r--r--	root/root	usr/man1/hcitool.1
-rw-r--r--	root/root	usr/man1/hidd.1
-rw-r--r--	root/root	usr/man1/l2ping.1
-rw-r--r--	root/root	usr/man1/pand.1
-rw-r--r--	root/root	usr/man1/rfcomm.1
-rw-r--r--	root/root	usr/man1/sdptool.1
drwxr-xr-x	root/root	usr/man8/
-rw-r--r--	root/root	usr/man8/bluetoothd.8
-rw-r--r--	root/root	usr/man8/hciattach.8
-rw-r--r--	root/root	usr/man8/hciconfig.8
drwxr-xr-x	root/root	usr/sbin/
-rwxr-xr-x	root/root	usr/sbin/bluetoothd
-rwxr-xr-x	root/root	usr/sbin/hciattach
-rwxr-xr-x	root/root	usr/sbin/hciconfig
drwxr-xr-x	root/root	usr/var/
drwxr-xr-x	root/root	usr/var/lib/
drwxr-xr-x	root/root	usr/var/lib/bluetooth/
